# rp/rpi-server.service
[Unit]
Description=rpi-server

[Service]
ExecStartPre=/bin/mkdir -p ~/home/kuka/LU-Recycling-Robot/local
ExecStartPre=/usr/bin/rsync -av --checksum --delete ~/home/kuka/LU-Recycling-Robot/rp/ /home/kuka/LU-Recycling-Robot/local/
ExecStart=/usr/bin/python -u "~/home/kuka/LU-Recycling-Robot/local/server.py"

WorkingDirectory=~/home/kuka/LU-Recycling-Robot/local
Restart=always

[Install]
WantedBy=default.target
